# mc-access-web.service
[Unit]
Description=Minecraft Access Log Web Viewer
After=network.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/mc-access-web.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
